# dist/aur/assistd-git/assistd.install
post_install() {
  cat <<'EOF'

==> Initialize your user config:
      assistd init-config
    (writes ~/.config/assistd/config.toml from defaults)

    Or copy the annotated system template manually:
      mkdir -p ~/.config/assistd
      cp /etc/assistd/config.toml ~/.config/assistd/config.toml

==> Install llama.cpp for the LLM backend (provides llama-server):
      pacman -S llama.cpp        # or build llama.cpp-* from the AUR

==> Enable and start the user service:
      systemctl --user enable --now assistd

==> CUDA acceleration for whisper STT (NVIDIA, optional):
      Rebuild from source with: cargo build --release --features cuda
      Requires the 'cuda' package and an NVIDIA GPU.

EOF
}

post_upgrade() {
  post_install
}
